# server/labbot.service
[Unit]
Description=Labbot server
After=mutli-user.target

[Service]
Type=simple
ExecStart=/home/meson800/labbot/server/env/bin/python -m labbot
WorkingDirectory=/home/meson800/labbot/server